# traits/special_purpose/rickroller/advanced/files/etc/rc.local
ifconfig at0 192.168.23.1
ifconfig at0 up
