# traits/special_purpose/nfc_login/files/etc/rc.local
/usr/bin/login-13mhz &
/usr/bin/login-125khz &
